# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net-device-up
stop on shutdown

respawn

chdir /var/www/SITENAME/source/
exec ../virtualenv/bin/gunicorn --user=ashley --bind unix:/tmp/SITENAME.socket superlists.wsgi:application
